# litellm-ollama/root/etc/s6-overlay/s6-rc.d/svc-litellm/run
#!/usr/bin/bash
/app/litellm --port ${LITELLM_PORT:-8000} --host ${LITELLM_HOST:-0.0.0.0}
