# litellm-ollama/root/etc/s6-overlay/s6-rc.d/init-ollama/run
#!/usr/bin/bash
# shellcheck shell=bash
# TODO Autopull anything from OLLAMA_PULL
echo "Hi"
